# infra/systemd/personel-secret-rotation.service
# =============================================================================
# Personel Platform — Secret Rotation Service
# Faz 5 Madde 55 — automated rotation of static secrets
# (postgres role passwords, vault approle secret IDs, MinIO root, KC admin).
# Cadence + per-secret behaviour: infra/scripts/secret-inventory.yaml
# =============================================================================

[Unit]
Description=Personel Platform Secret Rotation
Documentation=file:///opt/personel/docs/operations/secret-rotation.md
After=personel-compose.service docker.service
Requires=docker.service
PartOf=personel.target

[Service]
Type=oneshot
User=root
Group=root

WorkingDirectory=/opt/personel/infra
EnvironmentFile=/opt/personel/infra/compose/.env
# AUDIT_TOKEN must be provisioned out-of-band (not in compose .env);
# operator drops it into /etc/personel/secrets/audit-emit.token and the
# unit picks it up via systemd EnvironmentFile.
EnvironmentFile=-/etc/personel/secrets/audit-emit.env

# Any due secret triggers rotation; non-due secrets are skipped.
ExecStart=/opt/personel/infra/scripts/rotate-secrets.sh

TimeoutStartSec=600
Restart=no

# Hardening
ProtectSystem=strict
ReadWritePaths=/etc/personel/secrets /var/lib/personel/secret-rotation /var/log/personel /opt/personel/infra/compose
PrivateTmp=true
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=personel-secret-rotation

[Install]
WantedBy=personel.target
