# infra/systemd/personel-secret-rotation.timer
# =============================================================================
# Personel Platform — Secret Rotation Timer
# Runs weekly on Sundays at 04:00 local time. The rotate-secrets.sh script
# itself decides which secrets are due based on per-secret cadence in
# infra/scripts/secret-inventory.yaml — non-due secrets are skipped, so a
# weekly timer with mixed cadences is intentional.
# =============================================================================

[Unit]
Description=Personel Platform Secret Rotation Timer
PartOf=personel.target

[Timer]
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=15min
Persistent=true
Unit=personel-secret-rotation.service

[Install]
WantedBy=timers.target
